# services/recommender/root/etc/systemd/system/ponderosa-recommender.service
[Unit]
Description=Ponderosa Recommender Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ponderosa/recommender
ExecStart=/opt/ponderosa/recommender/run.sh
Restart=on-failure
RestartSec=5

Environment=SURVEYS_DIR=/var/ponderosa/surveys/pending
Environment=PROCESSED_DIR=/var/ponderosa/surveys/processed
Environment=RECOMMENDER_PIPE=/var/ponderosa/pipes/recommender
Environment=PRESENTER_PIPE=/var/ponderosa/pipes/presenter

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ponderosa-recommender

[Install]
WantedBy=multi-user.target
